# canopy_server_comm.service
[Unit]
Description = making canopy server comm up
[Service]
ExecStart = /bin/sh -c "canopy_server_comm --addr=$(hostname -I | cut -d' ' -f1):8080"